# supervisor.ini
[program:fortnite]
command=/usr/bin/gunicorn -b :5000 -w 4 fortnite:app
directory=/app
user=fortnite
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true

[program:nginx]
command=/usr/sbin/nginx
startsecs=5
stdout_events_enabled = true
stderr_events_enabled = true